# modd/crontab
0       16      *       *       *       /usr/local/bin/flip-inventory