# services/popularTeams.service
[Unit]
Description= log number of chatters in our server in python via tornado
After=network.target

[Service]
Type=simple
ExecStart= /opt/modd/popularTeams.sh

Restart=always

[Install]
WantedBy=multi-user.target
